# deploy/systemd/go-mall.service
[Unit]
Description=go-mall backend
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/mall/backend
Environment=GIN_MODE=release
ExecStart=/opt/mall/backend/bin/go-mall-server
Restart=always
RestartSec=3
TimeoutStopSec=30
LimitNOFILE=65535
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
