# deploy/systemd/go-mall-reconcile-refunds.timer
[Unit]
Description=Run go-mall refund reconciliation every five minutes

[Timer]
OnCalendar=*:0/5
Persistent=true
RandomizedDelaySec=30
Unit=go-mall-reconcile-refunds.service

[Install]
WantedBy=timers.target
