# deploy/systemd/go-mall-cancel-expired-orders.service
[Unit]
Description=go-mall cancel expired orders
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/mall/backend
ExecStart=/opt/mall/backend/bin/go-mall-cancel-expired-orders
