# deploy/systemd/go-mall-complete-shipped-orders.service
[Unit]
Description=go-mall complete shipped orders
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/mall/backend
ExecStart=/opt/mall/backend/bin/go-mall-complete-shipped-orders
